# rootfs/etc/inittab
::sysinit:/bin/mount -t proc proc /proc
::askfirst:-/bin/sh
